# src/DOBOT/service/dobot.service
[Unit]
Description=Dobot Magacian Server
After=network.target
StartLimitIntervalSec=0

[Service]
Type=simple
Restart=always
RestartSec=1
User=dobot
WorkingDirectory=/home/dobot/DobotMagicianServer
ExecStart=/usr/bin/python3 DobotServer.py

[Install]
WantedBy=multi-user.target
